# scripts/com.brain-memory.sleep-pass.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.brain-memory.sleep-pass</string>

    <!--
        ProgramArguments: bash + the committed, secret-free launchd wrapper.
        Placeholder substituted by scripts/setup-dogfood.sh:
          __WRAPPER__ → absolute path to scripts/sleep-pass-launchd.sh
        The wrapper sources the gitignored env file (path from BRAIN_MEMORY_SLEEP_ENV
        below) for keys + per-role config, then execs node on the compiled sleep-pass CLI.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__WRAPPER__</string>
    </array>

    <!--
        EnvironmentVariables: only the env-file path. All other config — node binary,
        sleep-pass JS, DB path, per-role providers, and API keys — lives in the
        gitignored env file (chmod 600), generated by scripts/setup-dogfood.sh.
        __ENV_FILE__ → absolute path to the sleep env file (BRAIN_MEMORY_SLEEP_ENV),
                       default ~/.config/brain-memory/sleep.env.
    -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>BRAIN_MEMORY_SLEEP_ENV</key>
        <string>__ENV_FILE__</string>
    </dict>

    <!--
        WatchPaths (L8N-01): on-write trigger — fires the same wrapper whenever the
        dirty-sentinel file is touched by EpisodicStore.append() after a new episode
        write. This cuts consolidation latency from up to ~1h to ~seconds after a write.
        __DIRTY_SENTINEL__ → absolute path to the sentinel file, substituted by
                             scripts/setup-dogfood.sh. Default:
                             ~/.config/brain-memory/.episodes-dirty
        ThrottleInterval=120s coalesces write bursts — launchd will not re-fire the
        job more than once per 120s, guarding the Anthropic judge against rapid bursts.
    -->
    <key>WatchPaths</key>
    <array>
        <string>__DIRTY_SENTINEL__</string>
    </array>

    <!-- ThrottleInterval: debounce/rate-limit for WatchPaths fires (seconds). -->
    <key>ThrottleInterval</key>
    <integer>120</integer>

    <!--
        StartCalendarInterval: hourly floor — run once per hour (at :00) after any wake.
        Kept alongside WatchPaths so the sleep pass still runs even when no writes occur
        (e.g. strength decay, eviction). launchd catches missed runs after sleep.
        Minute=0 only (no Hour key) means: fire at the top of every hour.
    -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- RunAtLoad false: do not fire immediately on launchctl load -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Redirect stdout/stderr to log files — sleep-pass-cli also writes here -->
    <key>StandardOutPath</key>
    <string>/tmp/brain-memory-sleep.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/brain-memory-sleep-err.log</string>
</dict>
</plist>
